# supervisord.conf
[supervisord]
nodaemon=true

[program:scribe]
command=scribed /etc/scribe.conf
autorestart=true
